# sci-mathematics/rstudio/files/rstudio-server.initd
#!/sbin/openrc-run
# Copyright 1999-2022 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

pidfile="/var/run/rstudio-server.pid"
command="/usr/bin/rserver"
command_args="--config-file /etc/rstudio/rserver.conf"

depend() {
	use net
}

